# service/gps_launcher.service
[Unit]
Description=GPS Mobility Launcher
After=multi-user.target

[Service]
Type=simple
User=<your-username>
WorkingDirectory=/home/<your-username>/pi-gps-logger
ExecStart=/home/<your-username>/pi-gps-logger/.venv/bin/python /home/<your-username>/pi-gps-logger/mobility_launcher.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target